# verity_gate.service
[Unit]
Description=Verity Gate BPF Security Policy
RequiresMountsFor=/sys/fs/bpf
After=local-fs.target
DefaultDependencies=no
RefuseManualStop=yes

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'exec /usr/sbin/verity_gate_loader $(sed -n "s/.*roothash=\\([^ ]*\\).*/\\1/p" /proc/cmdline)'

RemainAfterExit=yes
ProtectSystem=full
ProtectHome=yes
CapabilityBoundingSet=CAP_BPF CAP_SYS_ADMIN

[Install]
WantedBy=sysinit.target